# deploy/career-scout-ai.timer
[Unit]
Description=Career Scout AI Daily Run at 02:00 Warsaw Time

[Timer]
OnCalendar=*-*-* 02:00:00 Europe/Warsaw
Persistent=true
Unit=career-scout-ai.service

[Install]
WantedBy=timers.target
